# vps/polymarket-daily.timer
[Unit]
Description=Polymarket Engine — Daily Review Timer

[Timer]
OnCalendar=*-*-* 00:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target